# etc/profile.d/snow.csh
#!/bin/csh
set SNOWROOT /sNow
setenv PATH $SNOWROOT/snow-tools/bin:$PATH
setenv MANPATH $SNOWROOT/snow-tools/man:$MANPATH
setenv PDSH_RCMD_TYPE "ssh"
