# etc/profile.d/snow-source.csh
#!/bin/csh

set SNOW_PATH /sNow
setenv PATH ${SNOW_PATH}/snow-tools/bin:$PATH

