# etc/profile.d/snow.csh
#!/bin/csh
setenv PATH $SNOW_BIN:$PATH
setenv MANPATH $SNOW_MAN:$MANPATH
setenv PDSH_RCMD_TYPE "ssh"
