# etc/profile.d/snow.sh
#!/bin/bash

SNOWROOT=/sNow
export PATH=$SNOWROOT/snow-tools/bin:$PATH
export PDSH_RCMD_TYPE=ssh
