# etc/profile.d/snow.sh
#!/bin/bash
SNOW_PATH=/sNow
export PATH=$SNOW_PATH/snow-tools/bin:$PATH
export MANPATH=$SNOW_PATH/snow-tools/man:$MANPATH
export PDSH_RCMD_TYPE=ssh
