# etc/profile.d/snow-source.csh
#!/bin/csh
setenv PATH ${SNOW_BIN}:$PATH
